# apps/inventory-agent/src/InvAgent.Installers/linux/debian/postinst
#!/bin/sh
# postinst script for meridian-agent
set -e

INSTALL_DIR="/opt/meridian-agent"
CONFIG_DIR="/etc/meridian-agent"
CONFIG_FILE="${CONFIG_DIR}/config.json"
SERVICE_NAME="meridian-agent"
AGENT_USER="meridian-agent"

# Create dedicated service user if it doesn't exist
if ! id -u "${AGENT_USER}" >/dev/null 2>&1; then
    useradd --system --no-create-home --shell /usr/sbin/nologin "${AGENT_USER}"
fi

# Create directories
mkdir -p "${INSTALL_DIR}"
mkdir -p "${CONFIG_DIR}"

# Copy binary and assets (package places them in /usr/share/meridian-agent/)
if [ -d "/usr/share/meridian-agent" ]; then
    cp -r /usr/share/meridian-agent/. "${INSTALL_DIR}/"
fi

# Make binary executable
chmod +x "${INSTALL_DIR}/meridian-agent" 2>/dev/null || true

# Create default config if it doesn't exist
if [ ! -f "${CONFIG_FILE}" ]; then
    cat > "${CONFIG_FILE}" <<'DEFAULTCONFIG'
{
  "AgentConfig": {
    "ServerUrl": "https://your-meridian-server.com",
    "EnrollmentToken": "",
    "PrivacyTier": "full",
    "HeartbeatIntervalSeconds": 300,
    "InventoryIntervalSeconds": 14400,
    "LocalWebUiPort": 8787,
    "LocalQueueMaxSizeMb": 100,
    "LogLevel": "Information"
  }
}
DEFAULTCONFIG
fi

# Set ownership
chown -R "${AGENT_USER}:${AGENT_USER}" "${INSTALL_DIR}" "${CONFIG_DIR}"
chmod 640 "${CONFIG_FILE}"

# Install systemd unit
if [ -f "/usr/lib/meridian-agent/meridian-agent.service" ]; then
    cp "/usr/lib/meridian-agent/meridian-agent.service" /etc/systemd/system/
fi

# Enable and start the service
if command -v systemctl >/dev/null 2>&1; then
    systemctl daemon-reload
    systemctl enable "${SERVICE_NAME}"
    systemctl start "${SERVICE_NAME}" || true
fi
